# tests/full-sys/client/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin serverhud --noclear %I $TERM
